# scripts/gcp/oracle-scheduler.service
[Unit]
Description=Project Oracle Scheduler
After=network.target

[Service]
Type=simple
User=oracle
WorkingDirectory=/opt/project-oracle/src
Environment=PYTHONPATH=/opt/project-oracle/src
ExecStart=/opt/project-oracle/.venv/bin/python3 scheduler.py --start --day-of-week 0 --hour 8 --minute 0
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target
